# package_template/DEBIAN/postrm
#!/bin/bash
echo "* uninstalling !CHAOS and service"
update-rc.d -f chaos.sh remove
