# package_template/DEBIAN/postinst
#!/bin/bash
echo "* installing !CHAOS and service"
update-rc.d chaos.sh start 90 5 stop 90 0 1 6

