# package_template/DEBIAN/postrm
echo "* uninstalling !CHAOS and service from $INSTDIR"
rm -rf $INSTDIR > /dev/null
rm -f /etc/init/chaos-us.conf > /dev/null
rm -f /etc/init/chaos-cuiserver.conf > /dev/null
rm -f /etc/init/chaos-wanproxy.conf > /dev/null
update-rc.d -f chaos.sh remove > /dev/null
rm -rf /usr/bin/UnitServer > /dev/null 
rm -f /etc/ld.so.conf.d/chaos.conf
